# modules/intel/templates/mic-backup.cron.erb
#!/bin/sh
USER=root
DEST=/mnt/backup/intel-mic

SUBJECT="** DOVECOT CONF BACKUP ERROR **"
EMAIL="mcaberletti@<%= @mail_server %>"
MESSAGE=/tmp/message

if [ ! -d $DEST ]; then
	echo "Backup fail" > $MESSAGE
	echo "No destination directory found" >> $MESSAGE
	/bin/mail -s "$SUBJECT" "$EMAIL"  < $MESSAGE
	exit 1
fi

/bin/tar czf $DEST/mic-conf_`hostname -s`_`date -I`.tar.gz /etc/sysconfig/mic/* /opt/intel/mic/* 

if [ $? -eq 0 ]; then
	echo "Backup end Successfully :)"
else
	echo "Backup FAIL :(" > $MESSAGE
	/bin/mail -s "$SUBJECT" "$EMAIL"  < $MESSAGE
fi

# Rimozione dump più vecchi di 7 giorni
/usr/bin/find $DEST -type f -mtime +7 -exec rm {} \;
